# scheduler/crontab
*/15 * * * * /bin/bash /app/scripts/zoom-jobs.sh
*/10 * * * * /bin/bash /app/scripts/export-jobs.sh
5 * * * * /bin/bash /app/scripts/export-cleanup.sh
